# heroforge.service
[Unit]
Description=HeroForge Network Triage Web Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Development/HeroForge
ExecStart=/root/Development/HeroForge/target/release/heroforge serve --bind 127.0.0.1:8080 --database sqlite:///root/Development/HeroForge/heroforge.db
Restart=always
RestartSec=10

# Environment
Environment=RUST_LOG=info

# Security (adjust as needed)
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=heroforge

[Install]
WantedBy=multi-user.target
